# services/create/run
#!/bin/bash
set -euxo pipefail
exec 2>&1
cd /opt/create
exec python -m create.worker -c dev-create.conf --debug
