# services/create/run
#!/bin/bash
set -exo pipefail
cd /opt/create
. venv/bin/activate
exec python -m create.worker
