# services/healthcheck/run
#!/bin/bash
set -euxo pipefail
exec 2>&1
cd /opt/create
while :; do
    python -m create.healthcheck
    sleep 60
done
